# conf/supervisord.conf
[unix_http_server]
file=/path/to/dir/logs/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///path/to/dir/logs/supervisor.sock

[supervisord]
pidfile = /tmp/supervisord.pid
logfile = /path/to/dir/logs/supervisord.log

[program:php-fpm]
autorestart=true
command=/usr/sbin/php5-fpm -F
